# update_db/crontab
0 12 * * * python -u /app/db_update.py 
